# supervisord.conf
[supervisord]
nodaemon=true

[program:hbase]
command=/opt/hbase/bin/hbase master start

[program:thrift2]
command=/opt/hbase/bin/hbase thrift2 start